# deploy/etc/s6-overlay/s6-rc.d/web/run
#!/command/with-contenv execlineb

cd /app
/opt/pysetup/.venv/bin/gunicorn wsgi:application --bind 0.0.0.0:8000 --workers 2
