# deploy/etc/s6-overlay/s6-rc.d/worker/run
#!/command/with-contenv execlineb

cd /app
/opt/pysetup/.venv/bin/python manage.py qcluster
