# deploy/etc/s6-overlay/s6-rc.d/server/run
#!/command/with-contenv execlineb

cd /app
foreground { redirfd -r 0 create_user.py python manage.py shell }
gunicorn wsgi:application --config="deploy/gunicorn.conf.py"
