# etc/s6-overlay/s6-rc.d/worker/run
#!/bin/sh

cd /app
/opt/pysetup/.venv/bin/python manage.py setup_periodic_tasks
/opt/pysetup/.venv/bin/python manage.py qcluster
